# deploy/flotilla-watch.service
[Unit]
Description=flotilla watch — inter-agent relay + XO heartbeat clock
Documentation=https://github.com/jim80net/flotilla
After=network-online.target
Wants=network-online.target
# Stop a bad-token hot-loop from hammering Discord's IDENTIFY rate limit: give
# up after 5 failed starts in 5 minutes (manual restart required thereafter).
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
# --- Edit these paths for your install (see docs/watch-runbook.md). ExecStart
# assumes the binary is at ~/go/bin/flotilla (the `go install` default). ---
EnvironmentFile=-%h/.config/flotilla/flotilla-secrets.env
# Defense-in-depth for cold-boot DNS: after a power-failure reboot, systemd-
# resolved (127.0.0.53) can still be unanswering ~6s into the boot, so the relay
# gateway's discord.com lookup fails. The daemon now treats a failed relay open as
# NON-FATAL (clock-only + background retry — see cmd/flotilla/relay.go), so this
# pre-flight is a belt-and-braces nicety, not a correctness requirement: wait up to
# ~60s for discord.com to resolve, but ALWAYS exit 0 (even on timeout) so the
# safety-critical clock starts regardless. Never block the clock on DNS.
ExecStartPre=/bin/sh -c 'for i in $(seq 1 30); do getent hosts discord.com >/dev/null 2>&1 && exit 0; sleep 2; done; exit 0'
ExecStart=%h/go/bin/flotilla watch \
  --roster %h/.config/flotilla/flotilla.json \
  --secrets %h/.config/flotilla/flotilla-secrets.env
Restart=on-failure
# Restart=on-failure + the StartLimit above still guard GENUINELY-fatal startup
# errors (a malformed roster, an unknown surface, a missing operator_user_id with
# a configured channel) — those are config/programming errors the operator must
# fix, so a hot-loop-then-give-up is the right behavior. A transient relay-open
# failure is NO LONGER fatal (the daemon degrades to clock-only and retries the
# gateway in the background), so the cold-boot crash-loop that hit the
# 2026-06-10 power-failure reboot can no longer exhaust the burst limit.
RestartSec=15
KillSignal=SIGTERM
TimeoutStopSec=15

[Install]
WantedBy=default.target
